# com.example.substack-archiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Copy to ~/Library/LaunchAgents/ and edit the three absolute paths below
  to point at your local clone of substack-archiver. launchd does not
  expand $HOME, so the paths must be literal.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.substack-archiver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/substack-archiver/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/substack-archiver/archiver.log</string>

    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/substack-archiver/archiver.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
